# vagrant/corvus-web.ini
[group:corvus-web]
programs=api,task,redis,cache

[group:corvus-agent]
programs=agent-api,agent-task

[program:api]
command=/home/vagrant/corvus-web-env/bin/corvus_web dev_serve -p 10924
autostart=true
autorestart=unexpected
startsecs=3
startretries=3
stopsignal=TERM
stopwaitsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/corvus-web/web.log
stderr_logfile=/var/log/corvus-web/web.log
environment=CORVUS_WEB_CONFIG=/home/vagrant/corvus-web/vagrant/corvus_dashboard_config.py

[program:task]
command=/home/vagrant/corvus-web-env/bin/corvus_web task_daemon
autostart=true
autorestart=unexpected
startsecs=3
startretries=3
stopsignal=TERM
stopwaitsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/corvus-web/task.log
stderr_logfile=/var/log/corvus-web/task.log
environment=CORVUS_WEB_CONFIG=/home/vagrant/corvus-web/vagrant/corvus_dashboard_config.py

[program:redis]
command=/opt/redis/bin/redis-server
autostart=true
autorestart=unexpected
startsecs=3
startretries=3
stopsignal=QUIT
stopwaitsecs=5
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/corvus-web/redis.log
stderr_logfile=/var/log/corvus-web/redis.log

[program:cache]
command=/home/vagrant/corvus-web-env/bin/corvus_web cache_daemon
autostart=true
autorestart=unexpected
startsecs=3
startretries=3
stopsignal=TERM
stopwaitsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/corvus-web/cache.log
stderr_logfile=/var/log/corvus-web/cache.log
environment=CORVUS_WEB_CONFIG=/home/vagrant/corvus-web/vagrant/corvus_dashboard_config.py

[program:agent-api]
command=/home/vagrant/corvus-web-env/bin/corvus_web agent -p 23411
autostart=true
autorestart=unexpected
startsecs=3
startretries=3
stopsignal=TERM
stopwaitsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/corvus-web/agent-api.log
stderr_logfile=/var/log/corvus-web/agent-api.log
environment=CORVUS_WEB_CONFIG=/home/vagrant/corvus-web/vagrant/corvus_agent_config.py

[program:agent-task]
command=/home/vagrant/corvus-web-env/bin/corvus_web task_daemon
autostart=true
autorestart=unexpected
startsecs=3
startretries=3
stopsignal=TERM
stopwaitsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/corvus-web/agent-task.log
stderr_logfile=/var/log/corvus-web/agent-task.log
environment=CORVUS_WEB_CONFIG=/home/vagrant/corvus-web/vagrant/corvus_agent_config.py
