# vagrant/redis-cluster.ini
[program:redis-cluster-{{ port }}]
command=/opt/redis/bin/redis-server /etc/redis/redis-cluster-{{ port }}.conf
autostart=true
autorestart=unexpected
startsecs=3
startretries=3
stopsignal=QUIT
stopwaitsecs=5
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=syslog
stderr_logfile=syslog
